# deploy/it-department-api.service
# Manual /opt layout. scripts/setup-ubuntu.sh generates a unit with WorkingDirectory = your clone and User = repo owner.
[Unit]
Description=IT Department portal REST API (Fastify + MariaDB)
After=network-online.target mysql.service mariadb.service
Wants=network-online.target mariadb.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/it-department/backend
Environment=NODE_ENV=production
EnvironmentFile=/etc/it-department/api.env
ExecStart=/usr/bin/node dist/index.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
